# etc/cron.d/earnings-calendar
# 決算予定（決算発表スケジュール）日次同期
# 毎日 06:30（サーバーのローカルタイム）に実行
# 決算予定API（EDINET DB /v1/calendar）から当日〜90日後の決算予定（確定分）を取得し、
# 過去13ヶ月の実績履歴から各社の会計年度末・個社オフセットを逆算して予想分を補完、
# 発表が近い記録銘柄だけ個社API（決算短信の実績）で答え合わせした上で、
# EarningsSchedule を洗い替える。各日記の next_earnings_date を事前計算し、
# 決算前日（翌日決算）の通知をファンアウトする。
#
# API利用は1回の実行で 確定窓(〜3) + 履歴取得(〜15) + 個社答え合わせ(既定60)
# ＝ 80リクエスト前後（無料枠100リクエスト/日）。
# 1日2回以上に増やすと無料枠を超える可能性が高いため、cron は1日1回のまま維持する
# こと。個社答え合わせの予算を減らしたい場合は --verify-budget で調整できる
# （例: --verify-budget 40 で 確定+履歴+答え合わせ合計を60リクエスト前後に抑える）。
#
# 設定手順:
#   sudo cp /var/www/django/stock-dialy/etc/cron.d/earnings-calendar /etc/cron.d/
#   sudo chmod 644 /etc/cron.d/earnings-calendar
#
# ログ確認:
#   tail -f /var/www/django/stock-dialy/logs/earnings_calendar.log

# 分  時  日  月  曜日  ユーザー  コマンド
30   6   *   *   *     naoki    cd /var/www/django/stock-dialy && /var/www/django/stock-dialy/venv/bin/python manage.py sync_earnings_calendar >> /var/www/django/stock-dialy/logs/earnings_calendar.log 2>&1
