# etc/cron.d/earnings-calendar
# 決算予定（決算発表スケジュール）日次同期
# 毎日 06:30（サーバーのローカルタイム）に実行
# 決算予定API（EDINET DB /v1/calendar）から当日〜90日後の決算予定（確定分のみ）を
# 取得して EarningsSchedule を洗い替える。四半期末＋オフセット等の自前算出による
# 予定日は設定しない（実際の発表日とズレが大きいため廃止）。決算前日（翌日決算）の
# 通知もファンアウトする。
#
# API利用は1回の実行で確定窓の取得のみ（数リクエスト程度、無料枠100リクエスト/日）。
#
# 設定手順:
#   sudo cp /var/www/django/stock-dialy/etc/cron.d/earnings-calendar /etc/cron.d/
#   sudo chmod 644 /etc/cron.d/earnings-calendar
#
# ログ確認:
#   tail -f /var/www/django/stock-dialy/logs/earnings_calendar.log

# 分  時  日  月  曜日  ユーザー  コマンド
30   6   *   *   *     naoki    cd /var/www/django/stock-dialy && /var/www/django/stock-dialy/venv/bin/python manage.py sync_earnings_calendar >> /var/www/django/stock-dialy/logs/earnings_calendar.log 2>&1
